# deploy/launchd/com.forgefleet.forgefleetd.template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ForgeFleet daemon (forgefleetd) LaunchAgent template.

  Installed by `bootstrap-node-template.sh` on macOS hosts. Tokens
  `__USER_HOME__` and `__NODE_NAME__` are substituted at install time
  via sed; the resulting file lives at
  ~/Library/LaunchAgents/com.forgefleet.forgefleetd.plist and is
  loaded with `launchctl bootstrap gui/$(id -u) <path>` (which
  registers it for live `launchctl kickstart -k` restarts later, the
  happy path the wave dispatcher's Phase-2 expects).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.forgefleet.forgefleetd</string>

    <key>ProgramArguments</key>
    <array>
        <string>__USER_HOME__/.local/bin/forgefleetd</string>
        <string>--node-name</string>
        <string>__NODE_NAME__</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__USER_HOME__</string>
        <key>PATH</key>
        <string>__USER_HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/usr/sbin:/bin:/sbin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__USER_HOME__</string>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/.forgefleet/logs/forgefleetd.out.log</string>

    <key>StandardErrorPath</key>
    <string>__USER_HOME__/.forgefleet/logs/forgefleetd.err.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
